# systemd/astsync.service
[Unit]
Description=Sync Asterisk config files

[Service]
Type=oneshot
User=asterisk
Group=asterisk
ExecStart=/usr/local/sbin/astsync
